# com.arxiv.daily-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.arxiv.daily-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/envs/JLtorch/bin/python3</string>
        <string>/Users/Vir-G/Downloads/Projects/QLab-stack/Arxiv_filter.py</string>
        <string>--send</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Vir-G/Downloads/Projects/QLab-stack</string>

    <!-- 每天 15:00 触发；若此时休眠，唤醒后自动补执行 与美国时间对齐-->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>15</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- 开机 / 登录时也跑一次，防止休眠期间错过 -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/Vir-G/Downloads/Projects/QLab-stack/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Vir-G/Downloads/Projects/QLab-stack/launchd_stderr.log</string>
</dict>
</plist>
